# petalinux/xilinx-vck190-base-trd/project-spec/meta-base-trd/recipes-apps/trd-files/files/trd-autostart.service
[Unit]
Description=trd-autostart

[Service]
Type=oneshot
ExecStart=/usr/bin/trd-autostart.sh

[Install]
WantedBy=multi-user.target
